# ops/server-update/remoteserver-rollback.service
[Unit]
Description=RemoteServer rollback to the last backup

[Service]
Type=oneshot
ExecStart=/opt/remoteserver-update/rollback.sh
