# ops/server-update/remoteserver-rollback.path
[Unit]
Description=Watch for the RemoteServer rollback trigger

[Path]
PathExists=/var/lib/remoteserver/updates/rollback.trigger
Unit=remoteserver-rollback.service

[Install]
WantedBy=multi-user.target
